# scripts/hook_supervisor.conf
[program:webhook_coding_net]
autostart = true
user=your_user
command=/path/to/your/virtualenv/bin/uwsgi /path/to/your/hook_uwsgi.ini
priority=1
redirect_stderr=true
stopsignal=INT